# crontab
*/5 * * * * python manage.py refreshfeeds