# crontab
*/5 * * * * python /app/manage.py refreshfeeds